# systemd/gp-projects-watcher.service
[Unit]
Description=GitPal projects directory watcher
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/gp-projects-watcher
Restart=on-failure
RestartSec=10
Environment=HOME=%h
StandardOutput=null
StandardError=journal

[Install]
WantedBy=default.target
