# systemd/gp-server.service
[Unit]
Description=GitPal web dashboard (port 4242)
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/gp-server
Restart=on-failure
RestartSec=10
Environment=HOME=%h PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
StandardOutput=null
StandardError=journal

[Install]
WantedBy=default.target
