# systemd/gp-watcher@.service
[Unit]
Description=GitPal file watcher for %i
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/gp-watcher %h/projects/%i
Restart=on-failure
RestartSec=10
Environment=HOME=%h
StandardOutput=null
StandardError=journal

[Install]
WantedBy=default.target
